# infra/templates/allen.service
[Unit]
Description=Allen Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/allen/packages/server
EnvironmentFile=/home/ubuntu/allen/.env.production

# Cap V8 old-space heap so Node throws OOMError instead of growing past
# system RAM. 10 GiB leaves room for native allocations + Claude CLI +
# spawned MCP server subprocesses + OS buffers on a 16 GiB box.
Environment=NODE_OPTIONS=--max-old-space-size=10240

ExecStart=/usr/bin/node dist/app.js
StandardOutput=append:/var/log/allen/server.log
StandardError=append:/var/log/allen/server.log

# cgroup v2 memory limits — confines OOM to this service tree so the
# kernel's global OOM killer never touches sshd / SSM / network stack
# (which is what was leaving the instance unreachable until reboot).
MemoryAccounting=true
MemoryHigh=12G
MemoryMax=14G
# A subprocess (Claude CLI, MCP server, node-pty) being OOM-killed should
# not take down the whole allen unit — let the main process keep serving.
OOMPolicy=continue

# Cap subprocess + fd fanout. Each agent session can spawn 1 Claude CLI +
# up to 9 MCP servers; without limits a stuck session can fork unbounded.
TasksMax=16384
LimitNOFILE=65536

# Restart with backoff so a crash loop can't immediately re-OOM the box.
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
